# roles/wsl/files/etc/systemd/system/systemd-timesyncd.service.d/override.conf
# Managed by Ansible

[Unit]
ConditionVirtualization=
ConditionVirtualization=wsl
